# services/tts-service.service
[Unit]
Description=VoiceBridge TTS Service - 独立语音合成服务
After=network.target llm-service.service
Before=voicebridge.service

[Service]
Type=simple
User=gaofei
Group=gaofei
WorkingDirectory=/data/gaofei/00-code/VoiceBridge
ExecStart=/data/gaofei/00-code/VoiceBridge/bash/start_tts_systemd.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=voicebridge-tts

# 环境变量
Environment="PYTHONPATH=/data/gaofei/00-code/VoiceBridge"
Environment="TTS_HOST=0.0.0.0"
Environment="TTS_PORT=8011"

# 资源限制
LimitNOFILE=65535
LimitNPROC=4096

[Install]
WantedBy=multi-user.target